# co2-sensor.service
[Unit]
Description=CO2 Sensor API Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/share/co2-sensor
ExecStart=/usr/bin/python3 /usr/share/co2-sensor/co2_api.py --config=/etc/co2-sensor/config.yml
Restart=always
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
