# co2-sensor.service
[Unit]
Description=CO2 Sensor API Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/share/co2-sensor
ExecStart=/usr/bin/python3 -m flask --app co2_api run --host=0.0.0.0 --port=8080
Restart=always
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
